# deploy/systemd/pewcorder.service
[Unit]
Description=Pewcorder Django API (gunicorn)
After=network.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=pewcorder
Group=pewcorder
WorkingDirectory=/var/www/api.pewcorder.benlocher.com/backend
Environment=PYTHONPATH=/var/www/api.pewcorder.benlocher.com/backend
Environment=DJANGO_SETTINGS_MODULE=config.settings
EnvironmentFile=-/var/www/api.pewcorder.benlocher.com/backend/.env
ExecStart=/var/www/api.pewcorder.benlocher.com/.venv/bin/gunicorn \
  config.wsgi:application \
  --bind 127.0.0.1:8003 \
  --workers 2 \
  --timeout 300 \
  --access-logfile /home/pewcorder/logs/gunicorn-access.log \
  --error-logfile /home/pewcorder/logs/gunicorn-error.log
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
