# deploy/systemd/pewcorder-beat.service
[Unit]
Description=Pewcorder Celery beat
After=network.target redis-server.service postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=pewcorder
Group=pewcorder
WorkingDirectory=/var/www/api.pewcorder.benlocher.com/backend
Environment=PYTHONPATH=/var/www/api.pewcorder.benlocher.com/backend
Environment=DJANGO_SETTINGS_MODULE=config.settings
EnvironmentFile=-/var/www/api.pewcorder.benlocher.com/backend/.env
ExecStart=/var/www/api.pewcorder.benlocher.com/backend/.venv/bin/celery \
  -A config beat \
  --loglevel=INFO \
  --schedule=/var/www/api.pewcorder.benlocher.com/backend/celerybeat-schedule
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
